# systemd/gateway-ui.service
[Unit]
Description=Pi Gateway Web UI
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/opt/antigateway/app
ExecStart=/usr/bin/python3 /opt/antigateway/app/app.py
Restart=on-failure
RestartSec=5
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
